# localagentedit.service
[Unit]
Description=LocalAgent 24/7 Background Service
Documentation=file:///home/killer123/Desktop/vpn/README.md
After=network.target

[Service]
Type=simple
User=killer123
WorkingDirectory=/home/killer123/Desktop/vpn
ExecStart=/usr/bin/python3 /home/killer123/Desktop/vpn/agent_daemon.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
