# localagentedit.service
[Unit]
Description=LocalAgent 24/7 Background Service
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
User=%i
WorkingDirectory=%h/Desktop/vpn
ExecStart=/usr/bin/python3 %h/Desktop/vpn/agent_daemon.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
